# infra/latitude/boringd.service
[Unit]
Description=boringd — Firecracker microVM sandbox daemon (boring computers)
Documentation=https://github.com/boring-computers
After=network.target
Wants=network.target

[Service]
Type=simple
# Config (BORING_TOKEN, BORING_MAX, ...) is read from this env file.
# Leading '-' makes it optional so boringd still starts with no token set.
EnvironmentFile=-/etc/boring/boringd.env
ExecStart=/usr/local/bin/boringd
WorkingDirectory=/opt/boring
Restart=always
RestartSec=2
# boringd launches firecracker children which need /dev/kvm and unix sockets.
# It runs as root for the prototype (jailer/unprivileged is a hardening TODO).
User=root
Group=root
# Give firecracker room for per-VM sockets/overlays and pipes.
LimitNOFILE=1048576
LimitNPROC=65536
# Keep the working state on the box.
RuntimeDirectory=boring
StateDirectory=boring

[Install]
WantedBy=multi-user.target
